# files/sv/nginx/run
#!/bin/sh
exec 2>&1
exec /opt/kibana3/embedded/sbin/nginx -c /opt/kibana3/embedded/conf/nginx.conf

